# configuration/init/systemd/tedge-mapper-c8y.target
[Unit]
Description=tedge-mapper-c8y cloud profile services

[Install]
WantedBy=multi-user.target
