# configuration/init/systemd/tedge-cert-renewer.target
[Unit]
Description=thin-edge.io certificate renewal
