# configuration/init/systemd/tedge-mapper-az.target
[Unit]
Description=tedge-mapper-az cloud profile services

[Install]
WantedBy=multi-user.target
