# configuration/init/systemd/tedge-cert-renewer-c8y.target
[Unit]
Description=thin-edge.io certificate renewal for c8y cloud profiles
