# configuration/init/systemd/tedge-mapper-aws.target
[Unit]
Description=tedge-mapper-aws cloud profile services

[Install]
WantedBy=multi-user.target
