# tests/RobotFramework/tests/cumulocity/restart/on_startup.service
[Unit]
Description=Boot notifier

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/on_startup.sh

[Install]
WantedBy=multi-user.target