# configuration/package_scripts/_generated/tedge-flows/deb/postrm
#!/bin/sh
set -e

# Automatically added by thin-edge.io
if [ -d /run/systemd/system ]; then
	systemctl --system daemon-reload >/dev/null || true
fi
# End automatically added section
# Automatically added by thin-edge.io
if [ "$1" = "remove" ]; then
	if command -v deb-systemd-helper >/dev/null 2>&1; then
		deb-systemd-helper mask tedge-flows.service >/dev/null || true
	elif command -v systemctl >/dev/null 2>&1; then
		systemctl mask tedge-flows.service >/dev/null || true
	fi
fi

if [ "$1" = "purge" ]; then
	if command -v deb-systemd-helper >/dev/null 2>&1; then
		deb-systemd-helper purge tedge-flows.service >/dev/null || true
		deb-systemd-helper unmask tedge-flows.service >/dev/null || true
	elif command -v systemctl >/dev/null 2>&1; then
		systemctl unmask tedge-flows.service >/dev/null || true
	fi
fi
# End automatically added section
